# contrib/80-ndctl.rules
# do not edit this file, it will be overwritten on update

ACTION=="add", KERNEL=="nmem*", RUN+="ndctl-udev $kernel"
